# ansible/roles/vagrant-boxes/files/rc.local
#!/usr/bin/env bash

sleep 10
systemctl daemon-reload

exit 0
